# sys-apps/my-systemd-units/files/services/disable_wifi_powersave@.service
[Unit]
Description=Disable power control for wifi devices.

[Service]
ExecStart=/usr/sbin/iw dev %i set power_save off

[Install]
WantedBy=multi-user.target 
