# sys-apps/my-systemd-units/files/services/ntpd.service
[Unit]
Description=ntpd

[Service]
Restart=always
ExecStart=/usr/sbin/ntpd -d

[Install]
WantedBy=multi-user.target
