# sys-apps/my-systemd-units/files/services/plymouth-start.service
[Unit]
Description=Plymouth stub

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes
