# sys-apps/my-systemd-units/files/targets/ext_lan.target
[Unit]
Description=Extrenal lan devices
