# sys-apps/my-systemd-units/files/targets/hostapd.target
[Unit]
Description=Wireless hostapd ap's
