# sys-apps/my-systemd-units/files/targets/promiscuous.target
[Unit]
Description=Set wifi cards in promiscuous mode
