# sys-apps/my-systemd-units/files/auditd.service
[Unit]
Description=Auditd stub

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

