# sys-apps/my-systemd-units/files/distccd.socket
[Unit]
Description=Distccd Socket
Conflicts=distccd.service

[Socket]
ListenStream=3632
Accept=yes

[Install]
WantedBy=sockets.target
