# sys-apps/my-systemd-units/files/rtorrent.path
[Unit]
Description=rtorrent watcher

[Path]
PathChanged=/home/myuser/torrents

[Install]
WantedBy=multi-user.target
