# sys-apps/my-systemd-units/files/vixie-cron.service
[Unit]
Description=vixie-cron
After=local-fs.target syslog.service
Requires=local-fs.target syslog.service

[Service]
Type=forking
ExecStart=/usr/sbin/cron
PIDFile=/run/cron.pid

[Install]
Alias=cron.service
WantedBy=multi-user.target
